# core-agent/archiso/package/aegis-kiosk.service
[Unit]
Description=Aegis sphere kiosk (cage + Electron frontend)
Documentation=https://github.com/kevmill29/AegisOS
# Take over tty1 from the login getty and wait for seat management + devices.
Wants=seatd.service
After=seatd.service systemd-user-sessions.service systemd-logind.service
Conflicts=getty@tty1.service
After=getty@tty1.service
# Never give up on the kiosk — keep restarting it if it dies.
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/usr/bin/aegis-session
Restart=on-failure
RestartSec=3
TTYPath=/dev/tty1
StandardInput=tty
StandardOutput=journal
StandardError=journal
Environment=XDG_RUNTIME_DIR=/run/aegis

[Install]
WantedBy=multi-user.target
