# core-agent/archiso/airootfs/etc/systemd/system/aegis-install.service
[Unit]
Description=Aegis guided installer (live ISO)
# Pulled in ONLY by aegis-install.target (the text-installer boot entry sets
# systemd.unit=aegis-install.target). Never enable this unit in a regular
# target: its Conflicts= must not exist in a normal live boot's transaction.
After=getty@tty1.service systemd-user-sessions.service network-online.target
Wants=network-online.target
Conflicts=getty@tty1.service greetd.service

[Service]
Type=idle
ExecStart=/usr/local/bin/aegis-installer
# Give the installer a real controlling terminal on tty1 for archinstall's TUI.
StandardInput=tty
StandardOutput=tty
# stderr must reach the user too: bash's `read -p` prompt and archinstall's
# error text both go to stderr — with =journal the "Press Enter" prompt was
# invisible and the installer looked frozen after the banner.
StandardError=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
Environment=TERM=linux
Restart=no
