# core-agent/archiso/airootfs/etc/systemd/system/aegis-install.service
[Unit]
Description=Aegis guided installer (live ISO)
# Only when the "Install Aegis OS" boot entry was chosen (adds aegis.install to
# the kernel command line). A normal live boot skips this and greetd runs the
# sphere instead.
ConditionKernelCommandLine=aegis.install
After=getty@tty1.service systemd-user-sessions.service network-online.target
Wants=network-online.target
Conflicts=getty@tty1.service greetd.service

[Service]
Type=idle
ExecStart=/usr/local/bin/aegis-installer
# Give the installer a real controlling terminal on tty1 for archinstall's TUI.
StandardInput=tty
StandardOutput=tty
StandardError=journal
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
Environment=TERM=linux
Restart=no

[Install]
WantedBy=multi-user.target
