# core-agent/archiso/airootfs/etc/systemd/system/greetd.service.d/aegis-install-mode.conf
# On the live ISO, don't bring up the sphere kiosk when the user chose the
# "Install Aegis OS" boot entry (aegis.install on the kernel command line) —
# aegis-install.service takes tty1 for the installer instead. On an installed
# system aegis.install is never present, so this condition is always true and
# greetd starts normally.
[Unit]
ConditionKernelCommandLine=!aegis.install
